# overlay/etc/cron/crontabs/root
*/10 * * * * logrotate /etc/logrotate.conf --state=/var/lib/logrotate/logrotate.status 2>&1 | /usr/bin/logger -t logrotate
