# overlay/etc/cron/crontabs/root
*/10 * * * * logrotate /etc/logrotate.conf --state=/opt/var/logrotate.status
