# sbt-open-source/board/common/rootfs_overlay/etc/network/if-up.d/chrony-online
#!/bin/sh

# Mark the network time sources as online.
# chrony will automatically update the system time
# and RTC when it gets the real time from the network.
chronyc online

